# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mnt-data-supervisor-homeassistant-zhaquirks.mount
[Unit]
Description=HassOS ZHA Quirks
Wants=mnt-data.mount
After=mnt-data.mount

[Mount]
What=/usr/lib/zhaquirks/
Where=/mnt/data/supervisor/homeassistant/zhaquirks
Type=none
Options=bind

[Install]
WantedBy=local-fs.target
